# crontab
0 4 */5 * * /root/backup.sh
0 4 */5 * * /root/rotate.sh
